# abe/sakuracoin.conf
[program:sakuracoin]
user=coind
autostart=true
autorestart=true
command=/usr/bin/sakuracoind -rpcuser=u -rpcpassword=p -datadir=/mnt/datadir/sakuracoin
