# abe/sayacoin.conf
[program:sayacoin]
user=coind
autostart=true
autorestart=true
command=/usr/bin/sayacoind -rpcuser=u -rpcpassword=p -datadir=/mnt/datadir/sayacoin
