# abe/sha1coin.conf
[program:sha1coin]
user=sha1coin
autostart=true
autorestart=true
command=/usr/bin/sha1coind -rpcuser=u -rpcpassword=p -datadir=/mnt/datadir/sha1coin
