# abe/sv-abe.conf
[program:abe]
user=root
autostart=true
autorestart=true
command=python -m Abe.abe --config /abe.conf
