# abe/sv-nomp.conf
[program:nomp]
user=root
autostart=true
autorestart=true
directory=/nomp
command=nodejs init.js
