# abe/yaycoin.conf
[program:yaycoin]
user=coind
autostart=true
autorestart=true
command=/usr/bin/yaycoind -rpcuser=u -rpcpassword=p -datadir=/mnt/datadir/yaycoin
