# dnsseed/sv-web.conf
[program:web]
user=root
autostart=true
autorestart=true
directory=dnsseed-ruby
command=ruby web.rb -o 0.0.0.0
